# scripts/systemd/runtime-units/rastro-irado-collector.service.in
[Unit]
Description=IRAI Collector (WSL, escopo WIN$N/WDO$N)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@RUNTIME_ROOT@
ExecStart=/usr/bin/env bash @RUNTIME_ROOT@/scripts/systemd/collector-wsl.sh
Restart=on-failure
RestartSec=5
Environment=IRAI_PYTHON=/mnt/c/Users/brenoperucchi/AppData/Local/Microsoft/WindowsApps/py.exe
Environment=IRAI_PYTHON_VERSION_FLAG=-3.12
Environment=IRAI_COLLECTOR_INTERVAL=60

[Install]
WantedBy=default.target
