# scripts/systemd/rastro-irado-gex.timer
[Unit]
Description=IRAI GEX worker — após primeira barra WIN confirmar a sessão efetiva

[Timer]
# O LIVE resolve source/effective pelo ledger WIN. Às 09:10 a barra M5 de
# 09:00 já deve estar fechada/coletada; em feriado não há barra e o worker
# falha fechado sem publicar snapshot.
OnCalendar=Mon..Fri *-*-* 09:10:00 America/Sao_Paulo
Persistent=true
Unit=rastro-irado-gex.service

[Install]
WantedBy=timers.target
